# metrics-collector/supervisord.conf
[supervisord]
 
nodaemon=true
 
[program:tomcat7]
 
command=/bin/bash -c "/var/lib/tomcat7/bin/startup.sh"
