# metrics-collector/supervisord.conf
[supervisord]
 
nodaemon=true
 
[program:tomcat7]
 
command=/bin/bash -c "/opt/tomcat7/bin/catalina.sh run"
